# examples/nginx/uwsgi.service
[Unit]
Description=uWSGI instance to serve acme2certifier

[Service]
RuntimeDirectory=uwsgi
ExecStart=uwsgi --ini acme2certifier.ini
Restart=always
Type=notify
NotifyAccess=all
WorkingDirectory=/opt/acme2certifier
User=nginx

[Install]
WantedBy=multi-user.target
